# docker/etc/crontabs/nginx
0 8 * * * cd /var/www/app/kanboard && php83 ./cli cronjob >/dev/null 2>&1
